# templates/autojump.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=autojump - purge old entries

[Service]
Type=oneshot
ExecStart=/usr/bin/zsh -ic 'autojump --purge'
PrivateTmp=true
CPUQuota=60%
