# templates/restic-clean.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Cleanup safe dir backups

[Timer]
OnBootSec=6h
OnUnitActiveSec=6h
RandomizedDelaySec=10
Persistent=true

[Install]
WantedBy=timers.target
