# templates/restic-backup.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Backup safe dir

[Service]
Type=oneshot
ExecStart={{ ansible_user_dir }}/.bin/rbackup create
PrivateTmp=true
CPUQuota=60%
